# systemd/chellow_upgrader.timer
[Unit]
Description=Run Chellow Upgrader every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
