# config/chellow.conf
description "Chellow"
author      "Tony Locke"

start on runlevel [2345]
stop on shutdown


script
    echo $$ > /var/run/chellow.pid
    exec /var/www/chellow/venv/bin/chellow_start.sh >> /var/log/chellow.log 2>&1
end script

pre-start script
    echo "[`date`] Chellow Starting" >> /var/log/chellow.log
end script


pre-stop script
    rm /var/run/chellow.pid
    echo "[`date`] Chellow Stopping" >> /var/log/chellow.log
end script
